# scripts/com.marc4data.worldcup.autopull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.marc4data.worldcup.autopull</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/marcalexander/projects/ai_orchestrator_claude/world_cup_soccer_2026/scripts/auto_pull.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/marcalexander/Library/Logs/worldcup_autopull.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/marcalexander/Library/Logs/worldcup_autopull.err.log</string>
</dict>
</plist>
